# openvpn/Configs/adblock.service
[Unit]
Description=Adblock domain update
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/adblock.sh
