# Download/Configs/jackett.service
[Unit]
Description=Jackett

[Service]
User=jackett
Group=jackett
ExecStart=/opt/jackett/jackett

[Install]
WantedBy=multi-user.target
