# Download/Configs/radarr.service
[Unit]
Description=Radarr

[Service]
User=radarr
Group=radarr
ExecStart=/usr/bin/mono /opt/radarr/Radarr.exe

[Install]
WantedBy=multi-user.target
